# deploy/autodeploy.service
[Unit]
Description=GitHub Webhook Auto-Deploy (Multi-Project)
After=network.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/docker/nyokiiaapp/deploy
Environment=WEBHOOK_SECRET=CHANGE_ME_TO_YOUR_SECRET
Environment=WEBHOOK_PORT=9000
Environment=CONFIG_FILE=/docker/nyokiiaapp/deploy/projects.json
ExecStart=/usr/bin/python3 /docker/nyokiiaapp/deploy/webhook.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
